# install/com.sentinelcore.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>             <string>com.sentinelcore</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/python3</string>
    <string>/opt/sentinelcore/sentinelcore.py</string>
  </array>
  <key>RunAtLoad</key>         <true/>
  <key>KeepAlive</key>         <true/>
  <key>WorkingDirectory</key>  <string>/opt/sentinelcore</string>
  <key>StandardOutPath</key>   <string>/tmp/sentinelcore.log</string>
  <key>StandardErrorPath</key> <string>/tmp/sentinelcore.err</string>
</dict>
</plist>
